# rootfs/etc/services.d/dashboard-streams/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Finish script for Dashboard Streams
# ==============================================================================

bashio::log.info "Dashboard Streams stopped"
