# rootfs/etc/services.d/dashboard-streams/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Dashboard Streams service
# ==============================================================================

bashio::log.info "Starting Dashboard Streams..."

# Activate virtual environment
source /opt/venv/bin/activate

# Run the server
exec python3 /usr/src/app/server.py
